# templates/netapp-smo.service.erb
# Managed by Puppet
#
[Unit]
Description=Snapchat for Oracle


[Service]
ExecStart=<%= @smo_root %>/smo/bin/smo_server start
ExecStop=<%= @smo_root %>/smo/bin/smo_server stop
Type=forking

[Install]
WantedBy=multi-user.target

